# usr/lib/systemd/system/user@.service.d/30-get-vol.conf
[Service]
ExecStopPre=/bin/sh -c 'sudo -u droidian XDG_RUNTIME_DIR=/run/user/32011 /usr/bin/getvol > /var/lib/restore-volume/volume'
