# debian/adaptation-oneplus3-configs.tristate.service
[Unit]
Description=tristate switch implementation
After=sm.puri.Phosh.service

[Service]
Type=simple
ExecStart=/usr/bin/tristate
Restart=always
