# debian/adaptation-oneplus3-configs.fpc1020.service
[Unit]
Description=fingerprint touch button service
After=sm.puri.Phosh.service

[Service]
Type=simple
ExecStart=/usr/bin/fpc1020
Restart=always
