# debian/adaptation-oneplus3-configs.vendor-overlay.service
[Unit]
Description=Vendor overlay service
Before=lxc@android.service
After=android-mount.service

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'mount -o bind /usr/lib/droid-vendor-overlay/etc/init/com.qualcomm.qti.wifidisplayhal@1.0-service.rc /vendor/etc/init/com.qualcomm.qti.wifidisplayhal@1.0-service.rc'
ExecStart=/bin/sh -c 'mount -o bind /usr/lib/droid-vendor-overlay/etc/init/vendor.lineage.trust@1.0-service.rc /vendor/etc/init/vendor.lineage.trust@1.0-service.rc'
ExecStart=/bin/sh -c 'mount -o bind /usr/lib/droid-vendor-overlay/etc/init/vendor.qti.hardware.alarm@1.0-service.rc /vendor/etc/init/vendor.qti.hardware.alarm@1.0-service.rc'
ExecStart=/bin/sh -c 'mount -o bind /usr/lib/droid-vendor-overlay/etc/init/vendor.qti.hardware.soter@1.0-service.rc /vendor/etc/init/vendor.qti.hardware.soter@1.0-service.rc'
ExecStart=/bin/sh -c 'mount -o bind /usr/lib/droid-vendor-overlay/etc/init/android.hardware.vr@1.0-service.rc /vendor/etc/init/android.hardware.vr@1.0-service.rc'
ExecStart=/bin/sh -c 'mount -o bind /usr/lib/droid-vendor-overlay/etc/init/android.hardware.drm@1.0-service.rc /vendor/etc/init/android.hardware.drm@1.0-service.rc'
ExecStart=/bin/sh -c 'mount -o bind /usr/lib/droid-vendor-overlay/etc/init/android.hardware.cas@1.0-service.rc /vendor/etc/init/android.hardware.cas@1.0-service.rc'
ExecStart=/bin/sh -c 'mount -o bind /usr/lib/droid-vendor-overlay/etc/init/android.hardware.drm@1.1-service.clearkey.rc /vendor/etc/init/android.hardware.drm@1.1-service.clearkey.rc'
ExecStart=/bin/sh -c 'mount -o bind /usr/lib/droid-vendor-overlay/etc/init/android.hardware.drm@1.1-service.widevine.rc /vendor/etc/init/android.hardware.drm@1.1-service.widevine.rc'
ExecStart=/bin/sh -c 'mount -o bind /usr/lib/droid-vendor-overlay/etc/init/vendor.nxp.hardware.nfc@1.1-service.rc /vendor/etc/init/vendor.nxp.hardware.nfc@1.1-service.rc'
ExecStart=/bin/sh -c 'mount -o bind /usr/lib/droid-vendor-overlay/etc/init/hw/init.oneplus3.rc /vendor/etc/init/hw/init.oneplus3.rc'

[Install]
WantedBy=multi-user.target
