# etc/systemd/user/pulseaudio.service.d/30-droid-config.conf
[Service]
ExecStart=
ExecStart=/usr/bin/pulseaudio --daemonize=no --log-target=journal -n --file=/etc/pulse/droid.pa
