# etc/systemd/system/pulseaudio.service.d/30-get-vol.conf
[Service]
ExecStop=/usr/bin/bash -c '/usr/bin/getvol > /var/lib/restore-volume/volume'
